# Appendix/AP-1/crawler_with_celery_sample.ini
[program:crawler_with_celery_sample_download]
directory=/Users/peketamin/python/crawler_with_celery_sample
command=/Users/peketamin/python/crawler_with_celery_sample/venv_crawler_with_celery_sample/bin/celery -A crawler_with_celery_sample worker -Q download -c 2 -l warning -n download@%%h
autostart=true
autorestart=true
startsecs=10
numprocs=1
stopwaitsecs=600
stopasgroup=false
killasgroup=false
user=peketamin
stdout_logfile=/Users/peketamin/python/crawler_with_celery_sample/worker_download_stdout.log
stderr_logfile=/Users/peketamin/python/crawler_with_celery_sample/worker_download_stderr.log
environment=PATH="/Users/peketamin/python/crawler_with_celery_sample/venv_crawler_with_celery_sample/bin:/usr/local/bin:%(ENV_PATH)s"

[program:crawler_with_celery_sample_media]
directory=/Users/peketamin/python/crawler_with_celery_sample
command=/Users/peketamin/python/crawler_with_celery_sample/venv_crawler_with_celery_sample/bin/celery -A crawler_with_celery_sample worker -Q media -c 2 -l warning -n media@%%h
autostart=true
autorestart=true
startsecs=10
numprocs=1
stopwaitsecs=600
stopasgroup=false
killasgroup=false
user=peketamin
stdout_logfile=/Users/peketamin/python/crawler_with_celery_sample/worker_media_stdout.log
stderr_logfile=/Users/peketamin/python/crawler_with_celery_sample/worker_media_stderr.log
environment=PATH="/Users/peketamin/python/crawler_with_celery_sample/venv_crawler_with_celery_sample/bin:/usr/local/bin:%(ENV_PATH)s"
